# config_files/blinky.service
[Unit]
Description=FlaschPlayer - LED Display Controller
After=network.target
Wants=network.target

[Service]
Type=simple
ExecStart=/path/to/flaschplayer/.venv/bin/python3 /path/to/flaschplayer/entry_point.py
WorkingDirectory=/path/to/flaschplayer
Restart=on-failure
RestartSec=5
StartLimitBurst=36000
User=root
Environment="WORK_DIR=/path/to/flaschplayer"
Environment="NEOPIXEL=True"
Environment="BOT_TOKEN=YOUR_BOT_TOKEN"
Environment="ROOT=YOUR_TELEGRAM_USER_ID"

[Install]
WantedBy=multi-user.target
